# etc/systemd/system/mpd.target
[Unit]
Description=Music Player Daemon (all instances)
After=multi-user.target

[Install]
WantedBy=multi-user.target